# fswatchd.service
[Unit]
Description = Watches for FS events


[Service]
ExecStart=fswatchd


[Install]
WantedBy=multi-user.target
